# roles/desktop/files/config/systemd/user/borgmatic.timer
[Unit]
Description=Run borgmatic backup

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=20m

[Install]
WantedBy=timers.target
